# lambda/vmdk-webhook-watcher.service
# Place this file at /etc/systemd/system/
[Unit]
Description=VMDK Webhook Watcher Service
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/vmdk_webhook_watcher.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target